# docker/cron/crontab
5 * * * * cd /var/www && sudo -u www-data /usr/bin/wget -O - -q -t 1 https://localhost/cron.php >> /var/log/cron.log
